# supervisord-apache2.conf
[program:apache2]
command=/start-apache2.sh
autostart=true
autorestart=true